# big-data/hadoop2/hadoop2/etc/profile.d/hadoop2-home.sh
hadoopVersion=2.3.0
hadoopHome="/opt/hadoop-$hadoopVersion"
export HADOOP_HOME=$hadoopHome
export PATH=$PATH:$HADOOP_HOME/bin:$HADOOP_HOME/sbin
