# big-data/cassandra-mgmt/cassandra-mgmt/etc/profile.d/cassandra-mgmt-home.sh
cassandraHome="/opt/cassandra-2.0.4"
export CASSANDRA_HOME=$cassandraHome
export PATH=$PATH:$CASSANDRA_HOME/bin
