# big-data/solr/solr/etc/profile.d/solr-home.sh
solrHome="/opt/solr-4.4.0"
export SOLR_HOME=$solrHome
export PATH=$PATH:$SOLR_HOME/bin
