# big-data/nutch/nutch/etc/profile.d/nutch-home.sh
nutchHome="/opt/nutch-1.5"
export NUTCH_HOME=$nutchHome
export PATH=$PATH:$NUTCH_HOME/bin
