# big-data/storm/storm/etc/profile.d/storm-home.sh
stormHome="/opt/storm-0.8.2"
stormConf="/etc/storm"
export STORM_HOME=$stormHome
export STORM_CONF_DIR=$stormConf
export PATH=$PATH:$STORM_HOME/bin
