# big-data/storm/storm/etc/profile.d/storm-home.sh
stormHome="/opt/storm-0.8.2"
export STORM_HOME=$stormHome
export PATH=$PATH:$STORM_HOME/bin
