# management/server/server-karaf/src/main/system/etc/init/management.conf
# Upstart job to start the SS Management server 

description "SS Management server"

start on startup

stop on runlevel [!2345]

exec /opt/subutai-mng/bin/start

respawn

