# big-data/oozie-server/oozie-server/etc/profile.d/oozie-server-home.sh
oozieHome="/opt/oozie-3.3.2"
export OOZIE_HOME=$oozieHome
export PATH=$PATH:$OOZIE_HOME/bin
