# big-data/hadoop/hadoop/etc/profile.d/hadoop-home.sh
export HADOOP_HOME="/opt/hadoop-1.2.1"
export HADOOP_CONF_DIR="/etc/hadoop"
export PATH=$PATH:$HADOOP_HOME/bin:$HADOOP_HOME/sbin
