# big-data/hadoop/hadoop/etc/profile.d/hadoop-home.sh
#!/bin/bash
set -e

export HADOOP_HOME="/opt/hadoop-1.2.1"
export HADOOP_CONF_DIR="/etc/hadoop"
