# big-data/hadoop/hadoop/etc/profile.d/hadoop-home.sh
hadoopHome="/opt/hadoop-1.2.1"
export HADOOP_HOME=$hadoopHome
export PATH=$PATH:$HADOOP_HOME/bin
