# big-data/hbase/hbase/etc/profile.d/hbase-home.sh
hbaseHome="/opt/hbase-0.94.16"
export HBASE_HOME=$hbaseHome
export PATH=$PATH:$HBASE_HOME/bin
