# big-data/hbase/hbase/etc/profile.d/hbase-home.sh
hbaseVersion=0.98.3
hbaseHome="/opt/hbase-$hbaseVersion"
hbaseConf="/etc/hbase/"
export HBASE_HOME=$hbaseHome
export PATH=$PATH:$HBASE_HOME/bin
