# big-data/pig/pig/etc/profile.d/pig-home.sh
pigHome="/opt/pig-0.12.0"
export PIG_HOME=$pigHome
export PATH=$PATH:$PIG_HOME/bin
