# host-scripts/debian/subutai-host-scripts/etc/profile.d/subutai-share.sh
subutaiShare="/usr/share/subutai-host-scripts"
export SUBUTAI_SHARE=$subutaiShare
export PATH=$PATH:$SUBUTAI_SHARE/bin
