# big-data/presto/presto/etc/profile.d/presto-home.sh
prestoHome="/opt/presto-server-0.69"
export PRESTO_HOME=$prestoHome
export PATH=$PATH:$PRESTO_HOME/bin
