# big-data/hive/hive/etc/profile.d/hive-home.sh
hiveHome="/opt/hive"
hiveConf="/etc/hive"
export HIVE_HOME=$hiveHome
export HIVE_CONF_DIR=$hiveConf
export PATH=$PATH:$HIVE_HOME/bin
