# big-data/hive/hive/etc/profile.d/hive-home.sh
hiveHome="/opt/hive-0.11.0-bin"
export HIVE_HOME=$hiveHome
export PATH=$PATH:$HIVE_HOME/bin
