# big-data/accumulo/accumulo/etc/profile.d/accumulo-home.sh
accumuloHome="/opt/accumulo-1.4.2"
export ACCUMULO_HOME=$accumuloHome
export PATH=$PATH:$ACCUMULO_HOME/bin
