# big-data/activemq/activemq/etc/profile.d/activemq-home.sh
ActivemqHome="/opt/apache-activemq-5.9.1"
export ACTIVEMQ_HOME=$ActivemqHome
export PATH=$PATH:$ACTIVEMQ_HOME/bin
