# big-data/flume/flume/etc/profile.d/flume-home.sh
flumeHome="/opt/flume-1.4.0"
export FLUME_HOME=$flumeHome
export PATH=$PATH:$FLUME_HOME/bin
