# big-data/mahout/mahout/etc/profile.d/mahout-home.sh
mahoutHome="/opt/mahout"
export MAHOUT_HOME=$mahoutHome
export PATH=$PATH:$MAHOUT_HOME/bin
