# big-data/derby/derby/etc/profile.d/derby-home.sh
derbyHome="/opt/db-derby-10.4.2.0-bin"
export DERBY_HOME=$derbyHome
export PATH=$PATH:$DERBY_HOME/bin
