# big-data/mongo/mongo/etc/profile.d/mongo-home.sh
mongoHome="/opt/mongo"
export MONGO_HOME=$mongoHome
