# big-data/oozie-client/oozie-client/etc/profile.d/oozie-client-home.sh
oozieClientHome="/opt/oozie-client-3.3.2"
export OOZIE_CLIENT_HOME=$oozieClientHome
export PATH=$PATH:$OOZIE_CLIENT_HOME/bin
