# big-data/zookeeper/zookeeper/etc/profile.d/zookeeper-home.sh
zookeeperHome="/opt/zookeeper-3.4.5"
zookeeperConf="/etc/zookeeper"
export ZOOKEEPER_HOME=$zookeeperHome
export ZOOKEEPER_CONF_DIR=$zookeeperConf
export PATH=$PATH:$ZOOKEEPER_HOME/bin
