# big-data/zookeeper/zookeeper/etc/profile.d/zookeeper-home.sh
zookeeperHome="/opt/zookeeper-3.4.5"
export ZOOKEEPER_HOME=$zookeeperHome
export PATH=$PATH:$ZOOKEEPER_HOME/bin
