# big-data/shark/shark/etc/profile.d/shark-home.sh
sharkHome="/opt/shark-0.9.1"
export SHARK_HOME=$sharkHome
export PATH=$PATH:$SHARK_HOME/bin
