# big-data/lucene/lucene/etc/profile.d/lucene-home.sh
luceneHome="/opt/lucene-4.4.0"
export LUCENE_HOME=$luceneHome
export PATH=$PATH:$LUCENE_HOME/bin
