# big-data/cassandra/cassandra/etc/profile.d/cassandra-home.sh
cassandraHome="/opt/cassandra-2.0.9"
export CASSANDRA_HOME=$cassandraHome
export PATH=$PATH:$CASSANDRA_HOME/bin
