# big-data/sqoop/sqoop/etc/profile.d/sqoop-home.sh
sqoopHome="/opt/sqoop-1.4.3"
hadoopHome="/opt/hadoop-1.2.1"
export SQOOP_HOME=$sqoopHome
export PATH=$PATH:$SQOOP_HOME/bin
export HADOOP_COMMON_HOME=$hadoopHome
export HADOOP_MAPRED_HOME=$hadoopHome

